# config/supervisor-celery.conf
[program:celery-workers]
command=celery -A app.celery_app worker --concurrency=1 --loglevel=info
directory=/root/.openclaw/workspace/vault/projects/orchestrator
numprocs=0  # Initial count: 0, no resource usage
autostart=false  # Do not auto-start
autorestart=true  # Auto-restart if crashed
startsecs=5  # Must run 5 seconds to be considered successful
stopwaitsecs=10  # Wait 10 seconds for graceful shutdown
priority=999  # Low priority
stdout_logfile=/root/.openclaw/workspace/vault/projects/orchestrator/logs/supervisor-celery.log
stderr_logfile=/root/.openclaw/workspace/vault/projects/orchestrator/logs/supervisor-celery.err
user=root
environment=PYTHONUNBUFFERED="1"
